# sys-apps/usbguard/files/runscript
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="Daemon to protect your computer against BadUSB."
PIDFILE=/var/run/usbguard.pid

depend() {
	need localmount udev dbus
	use logger
}

start() {
	ebegin "Starting USBGuard"
	start-stop-daemon --start --exec /usr/sbin/usbguard-daemon --pidfile "${PIDFILE}" --make-pidfile --background -- -s -C -p "${PIDFILE}"
	eend $? "Failed to start USBGuard"
}

stop() {
	ebegin "Stopping USBGuard"
	start-stop-daemon --stop --pidfile "${PIDFILE}" --exec /usr/sbin/usbguard-daemon
	eend $? "Failed to stop USBGuard"
}
